# n8n-assistant.service
[Unit]
Description=Personal Work Tracker Assistant
After=network.target

[Service]
Type=simple
User=shivang
WorkingDirectory=/home/shivang/n8n
ExecStart=/home/shivang/n8n/start.sh
Restart=always
RestartSec=10
EnvironmentFile=/home/shivang/n8n/.env

# Logging
StandardOutput=append:/home/shivang/n8n/logs/service.log
StandardError=append:/home/shivang/n8n/logs/service.log

[Install]
WantedBy=multi-user.target
